# systemd/ethminer_exporter.service
[Unit]
Description=EthMiner_Exporter

[Service]
ExecStart=/opt/prometheus/ethminer_exporter
Restart=always

[Install]
WantedBy=multi-user.target
