# tools/ethminer_watchdog/ethminer_watchdog.service
[Unit]
Description=Kill Ethminers after SIGSEGV
Wants=ethminer_watchdog.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ethminer_watchdog
